# .puppet/files/init.d-elkproxy
#!/bin/bash
# chkconfig:   2345 80 20
set -e

case "$1" in
    start|restart) o=' -b' ;;
    *) o='' ;;
    esac

/usr/bin/python2 -m libelkproxy.elkproxyd -p /var/run/elkproxyd.pid -c /etc/elkproxyd$o "$1"
